# development/openjdk/profile.d/openjdk.csh
#!/bin/csh
setenv OPENJDK_HOME /usr/lib/openjdk
setenv MANPATH ${MANPATH}:${OPENJDK_HOME}/man
setenv PATH ${PATH}:${OPENJDK_HOME}/bin:${OPENJDK_HOME}/jre/bin
